# modules/tile_gen/cron.d/ofm_tile_gen
# Define common variables
CMD="sudo /data/ofm/venv/bin/python -u /data/ofm/tile_gen/bin/tile_gen.py"
LOG_DIR=/data/ofm/tile_gen/logs

# every hour, make a monaco run
10 * * * * ofm    $CMD make-tiles monaco --upload >> $LOG_DIR/monaco-make-tiles.log 2>&1

# every minute, set monaco to latest
* * * * * ofm     $CMD set-version monaco >> $LOG_DIR/monaco-set-version.log 2>&1

# every Tuesday morning at 00:10, make a planet run
10 0 * * 2 ofm    $CMD make-tiles planet --upload >> $LOG_DIR/planet-make-tiles.log 2>&1

# every Saturday at 8:00, set planet to latest
0 9 * * 6 ofm     $CMD set-version monaco >> $LOG_DIR/monaco-set-version.log 2>&1

# once per minute, create indexes
* * * * * ofm     $CMD make-indexes >> $LOG_DIR/make-indexes-cron.log 2>&1
