# modules/tile_gen/cron.d/ofm_tile_gen
# Define common variables
CMD="sudo /data/ofm/venv/bin/python -u /data/ofm/tile_gen/bin/tile_gen.py"
LOG_DIR=/data/ofm/tile_gen/logs

# every day at 23:10, make a monaco run
10 23 * * * ofm    $CMD make-tiles monaco --upload >> $LOG_DIR/monaco-make-tiles.log 2>&1

# debug every 15 minutes
#*/15 * * * * ofm    $CMD make-tiles monaco --upload >> $LOG_DIR/monaco-make-tiles.log 2>&1

# every minute, set monaco to latest
* * * * * ofm     $CMD set-version monaco >> $LOG_DIR/monaco-set-version.log 2>&1

# every Wednesday, make a planet run
10 0 * * 3 ofm    $CMD make-tiles planet --upload >> $LOG_DIR/planet-make-tiles.log 2>&1

# every Saturday, set planet to latest
0 11 * * 6 ofm     $CMD set-version planet >> $LOG_DIR/planet-set-version.log 2>&1

# once per minute, create indexes
* * * * * ofm     $CMD make-indexes >> $LOG_DIR/make-indexes-cron.log 2>&1
