# scripts/http_host/cron.d/ofm_ledns_reader
# once per day
2 34 * * * ofm     sudo bash /data/ofm/http_host/bin/ledns_reader.sh >> /data/ofm/http_host/logs/ledns_reader.log 2>&1'


