# scripts/loadbalancer/cron.d/ofm_loadbalancer
# every minute
* * * * * ofm     sudo /data/ofm/venv/bin/python -u /data/ofm/loadbalancer/loadbalancer.py fix >> /data/ofm/loadbalancer/logs/check.log 2>&1


